# supervisord/sprinklerd.conf
[program:sprinklerd]
command=sprinklerd --config=/etc/sprinklerd.ini
numprocs=1
directory=/tmp
autostart=true
user=eli
redirect_stderr=true
