# supervisord/sprinklerd.conf
[program:sprinklerd]
command=/usr/bin/python /home/eli/src/pysprinklers/sprinklerd.py
numprocs=1
directory=/tmp
autostart=true
user=eli
